# system/calp.service
[Unit]
Description=Calendar server

[Service]
User=calp
ExecStart=calp --config /etc/calp/config.scm -- server
Restart=always
